# packaging/postinst
#!/bin/sh
set -e

PKG=deltafeedback
ETC=/etc/$PKG
RUN=/var/lib/$PKG

# Dedicated unprivileged user.
if ! getent group $PKG >/dev/null; then
    addgroup --system $PKG
fi
if ! getent passwd $PKG >/dev/null; then
    adduser --system --ingroup $PKG --home $RUN \
        --no-create-home --gecos 'DeltaFeedback' --shell /usr/sbin/nologin $PKG
fi

# State directory (SQLite db, DC db, blobs).
install -d -o $PKG -g $PKG -m 0750 $RUN

# Initialise config.ini from the example ONLY if it does not yet exist —
# upgrades MUST NOT overwrite a user's edited config.
if [ ! -e $ETC/config.ini ]; then
    install -m 0644 -o root -g root $ETC/config.example.ini $ETC/config.ini
    # Make the bundled defaults point at the package layout, and split mutable
    # state (addr/mail_pw/hmac_secret) into a writable account file owned by
    # the service user — /etc/ stays root-owned read-only.
    sed -i \
        -e "s|^db=.*|db=$RUN/deltafeedback.db|" \
        -e "s|^web_root=.*|web_root=/usr/share/$PKG/web|" \
        $ETC/config.ini
    if ! grep -q '^account_path=' $ETC/config.ini; then
        echo "account_path=$RUN/account.ini" >> $ETC/config.ini
    fi
    echo "$PKG: created $ETC/config.ini. Next steps:"
    echo "    sudo -u $PKG $PKG --register <chatmail-domain> $ETC/config.ini"
    echo "    sudo systemctl enable --now $PKG.service"
    echo "    sudo -u $PKG $PKG --invite $ETC/config.ini    # show invite URL"
fi

# Account file: the service writes addr/mail_pw/hmac_secret here at registration
# and on first run. Must be writable by the deltafeedback user.
if [ ! -e $RUN/account.ini ]; then
    install -m 0600 -o $PKG -g $PKG /dev/null $RUN/account.ini
fi

WEB=/usr/share/$PKG/web
DEFAULTS=/usr/share/$PKG/welcome.defaults

# Restore a backup made by preinst when upgrading FROM a version that
# shipped welcome.* in the package data (and which dpkg therefore wiped
# during post-unpack cleanup). Idempotent: no-op if no backup exists.
for L in ru en; do
    if [ -e "$WEB/welcome.$L.html.upgrade-saved" ]; then
        mv -f "$WEB/welcome.$L.html.upgrade-saved" "$WEB/welcome.$L.html"
    fi
done

# Seed welcome HTML in the served web dir from the bundled defaults — only
# if it isn't there already. The package never ships welcome.* INSIDE
# /usr/share/$PKG/web/, so dpkg can't clobber an admin-edited file on
# upgrade. (build-deb.sh moves them into /usr/share/$PKG/welcome.defaults/.)
for L in ru en; do
    if [ ! -e $WEB/welcome.$L.html ] && [ -e $DEFAULTS/welcome.$L.html ]; then
        install -m 0644 $DEFAULTS/welcome.$L.html $WEB/welcome.$L.html
    fi
done

# Always ensure the example file is readable for diff'ing on upgrade.
chmod 0644 $ETC/config.example.ini || true

if [ -d /run/systemd/system ]; then
    systemctl daemon-reload || true
fi

#DEBHELPER#
exit 0
